# tornado_celery_docker.ini
[program:tornado]
command=/app/rpc_server.py --logging=info --log-file-prefix=/var/log/rpc_server.log
directory=/app
user=root
numprocs=1
stdout_logfile=/app/logs/tornado.log
stderr_logfile=/app/logs/tornado.log
autostart=true
autorestart=true
startsecs=10
priority=997


[program:tornado_celery]
command=celery worker -A messaging.tasks -E -l=INFO
directory=/app
user=root
numprocs=1
stdout_logfile=/app/logs/celery-worker.log
stderr_logfile=/app/logs/celery-worker.log
autostart=true
autorestart=true
startsecs=10

stopwaitsecs = 600
killasgroup=true
priority=998


[program:flower]
command=celery flower
directory=/app
user=root
numprocs=1
stdout_logfile=/app/logs/flower.log
stderr_logfile=/app/logs/flower.log
autostart=true
autorestart=true
startsecs=10
priority=999

